# data/btrfsd.service.in
[Unit]
Description=Btrfs maintenance daemon
After=local-fs.target fstrim.service
Wants=local-fs.target

[Service]
Type=simple
ExecStart=@BTRFSD_INSTALL_BIN@
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
