# data/btrfsd.timer
[Unit]
Description=Perform Btrfs maintenance actions
Documentation=man:btrfsd
ConditionVirtualization=!container

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
